# etc/services/systemd/disco.service
[Unit]
Description=Disco - Discord Rich Presence client

[Service]
ExecStart=/path/to/disco --retry-after <int>

[Install]
WantedBy=default.target
